# systemd/pi-control.service
[Unit]
Description=Pi Control Dashboard
After=network-online.target
Wants=network-online.target

[Service]
WorkingDirectory=/home/admin/pi-control
ExecStart=/usr/bin/node /home/admin/pi-control/server.js
Restart=on-failure
RestartSec=5
User=admin
Environment=NODE_ENV=production
Environment=MAP_TILE_URL=http://127.0.0.1:8090/styles/bright/{z}/{x}/{y}.png
Environment=MAP_TILE_ATTRIB=(local Colorado tiles)
Environment=MAP_TILE_FALLBACK_URL=https://{s}.tile.openstreetmap.org/{z}/{x}/{y}.png
Environment=MAP_TILE_MAX_ZOOM=17
Environment=MAP_TILE_MAX_NATIVE_ZOOM=17
# If you use nvm or non-system node, adjust ExecStart accordingly.

[Install]
WantedBy=multi-user.target
